# coverage-matrix/current_coverage.json
{
  "generated": "2026-03-26",
  "scope": "Enterprise endpoint telemetry. Multi-source SIEM ingestion across EDR, identity, DNS, and cloud planes.",
  "siem_version": "Elastic Cloud 9.2.x",
  "edr_version": "CrowdStrike Falcon",
  "techniques": {
    "T1003.001": {
      "name": "OS Credential Dumping: LSASS Memory",
      "hypothesis": "DH-001",
      "elastic_score": 2,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "notes": "Rules written, lab validation pending",
      "rules": [
        "LSASS Memory Dump via Rundll32 comsvcs",
        "Suspicious LSASS Process Access via Sysmon"
      ],
      "rule_count": 2
    },
    "T1558.003": {
      "name": "Steal or Forge Kerberos Tickets: Kerberoasting",
      "hypothesis": "DH-002",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "Kerberoasting RC4 Ticket Request"
      ],
      "rule_count": 1
    },
    "T1087": {
      "name": "Account Discovery",
      "hypothesis": "DH-003",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "Active Directory Enumeration via LDAP and Net Commands"
      ],
      "rule_count": 1
    },
    "T1550.002": {
      "name": "Use Alternate Authentication Material: Pass the Hash",
      "hypothesis": "DH-004",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "Pass-the-Hash NTLM Lateral Movement"
      ],
      "rule_count": 1
    },
    "T1071.004": {
      "name": "Application Layer Protocol: DNS",
      "hypothesis": "DH-005, DH-019",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "notes": "DH-019 requires Cloudflare gateway DNS logs (operational). DH-005 requires lab DNS debug logging.",
      "rules": [
        "High-Entropy DNS Subdomain Query"
      ],
      "rule_count": 1
    },
    "T1218": {
      "name": "System Binary Proxy Execution",
      "hypothesis": "DH-006",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "LOLBin Execution with Suspicious Arguments"
      ],
      "rule_count": 1
    },
    "T1003.006": {
      "name": "OS Credential Dumping: DCSync",
      "hypothesis": "DH-007",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "DCSync Attack from Non-Domain Controller"
      ],
      "rule_count": 1
    },
    "T1053.005": {
      "name": "Scheduled Task/Job: Scheduled Task",
      "hypothesis": "DH-008",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "Suspicious Scheduled Task Creation"
      ],
      "rule_count": 1
    },
    "T1547.001": {
      "name": "Boot or Logon Autostart Execution: Registry Run Keys",
      "hypothesis": "DH-009",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "Registry Run Key Persistence"
      ],
      "rule_count": 1
    },
    "T1110.003": {
      "name": "Brute Force: Password Spraying",
      "hypothesis": "DH-010",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "Password Spray Attack Multiple Accounts"
      ],
      "rule_count": 1
    },
    "T1490": {
      "name": "Inhibit System Recovery",
      "hypothesis": "DH-012",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "Shadow Copy Deletion"
      ],
      "rule_count": 1
    },
    "T1621": {
      "name": "Multi-Factor Authentication Request Generation",
      "hypothesis": "DH-013",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "MFA Push Bombing via Duo"
      ],
      "rule_count": 1
    }
  },
  "summary": {
    "total_techniques": 12,
    "validated_elastic": 0,
    "validated_crowdstrike": 0,
    "backlog": 12,
    "score_3_rules": 0,
    "lab_status": "Detection lab build in progress. AWS isolated VPC architecture.",
    "score_2_rules": 1,
    "score_1_rules": 11,
    "score_0_no_coverage": 0,
    "elastic_coverage_pct": 0.0
  }
}