# coverage-matrix/current_coverage.json
{
  "generated": "2026-04-04",
  "scope": "Enterprise endpoint telemetry. Multi-source SIEM ingestion across EDR, identity, DNS, and cloud planes.",
  "siem_version": "Elastic Cloud 9.3.2",
  "edr_version": "CrowdStrike Falcon",
  "platform_architecture": "Two-repo model: detection-engineering-library (public, rules and methodology) + detection-platform (private, AlertTriageService and lab configs)",
  "techniques": {
    "T1003.001": {
      "name": "OS Credential Dumping: LSASS Memory",
      "hypothesis": "DH-001",
      "elastic_score": 2,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "notes": "Rules written, designated as POC for first ART validation loop. Two Sigma rules covering comsvcs.dll and Sysmon ProcessAccess paths.",
      "rules": [
        "Suspicious LSASS Process Access via Sysmon",
        "LSASS Memory Dump via Rundll32 comsvcs"
      ],
      "rule_count": 2
    },
    "T1558.003": {
      "name": "Steal or Forge Kerberos Tickets: Kerberoasting",
      "hypothesis": "DH-002",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "notes": "Rule written. Lab requires service account with SPN for realistic test.",
      "rules": [
        "Kerberoasting RC4 Ticket Request"
      ],
      "rule_count": 1
    },
    "T1087": {
      "name": "Account Discovery",
      "hypothesis": "DH-003",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "Active Directory Enumeration via LDAP and Net Commands"
      ],
      "rule_count": 1
    },
    "T1550.002": {
      "name": "Use Alternate Authentication Material: Pass the Hash",
      "hypothesis": "DH-004",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "Pass-the-Hash NTLM Lateral Movement"
      ],
      "rule_count": 1
    },
    "T1071.004": {
      "name": "Application Layer Protocol: DNS",
      "hypothesis": "DH-005, DH-019",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "notes": "DH-019 uses Cloudflare Gateway DNS logs (operational, dominant ingestion source). DH-005 requires lab DNS debug logging.",
      "rules": [
        "High-Entropy DNS Subdomain Query"
      ],
      "rule_count": 1
    },
    "T1218": {
      "name": "System Binary Proxy Execution",
      "hypothesis": "DH-006",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "notes": "Highest tuning complexity. FP check critical before production promotion.",
      "rules": [
        "LOLBin Execution with Suspicious Arguments"
      ],
      "rule_count": 1
    },
    "T1003.006": {
      "name": "OS Credential Dumping: DCSync",
      "hypothesis": "DH-007",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "notes": "High confidence signal. Non-DC source performing replication is near-certain malicious.",
      "rules": [
        "DCSync Attack from Non-Domain Controller"
      ],
      "rule_count": 1
    },
    "T1053.005": {
      "name": "Scheduled Task/Job: Scheduled Task",
      "hypothesis": "DH-008",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "Suspicious Scheduled Task Creation"
      ],
      "rule_count": 1
    },
    "T1547.001": {
      "name": "Boot or Logon Autostart Execution: Registry Run Keys",
      "hypothesis": "DH-009",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "rules": [
        "Registry Run Key Persistence"
      ],
      "rule_count": 1
    },
    "T1110.003": {
      "name": "Brute Force: Password Spraying",
      "hypothesis": "DH-010",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "notes": "Duo telemetry is primary source (operational). Entra SignInLogs blocked on P1/P2 license.",
      "rules": [
        "Password Spray Attack Multiple Accounts"
      ],
      "rule_count": 1
    },
    "T1490": {
      "name": "Inhibit System Recovery",
      "hypothesis": "DH-012",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "notes": "Near-zero FP rate. Shadow copy deletion has essentially no legitimate use case.",
      "rules": [
        "Shadow Copy Deletion"
      ],
      "rule_count": 1
    },
    "T1621": {
      "name": "Multi-Factor Authentication Request Generation",
      "hypothesis": "DH-013",
      "elastic_score": 1,
      "crowdstrike_score": 0,
      "status": "Backlog",
      "notes": "Duo Trust Monitor already active. This hypothesis formalizes the detection rule for correlation.",
      "rules": [
        "MFA Push Bombing via Duo"
      ],
      "rule_count": 1
    }
  },
  "summary": {
    "total_techniques": 12,
    "validated_elastic": 0,
    "validated_crowdstrike": 0,
    "backlog": 12,
    "score_3_rules": 0,
    "score_2_rules": 1,
    "score_1_rules": 11,
    "score_0_no_coverage": 0,
    "elastic_coverage_pct": 0.0,
    "lab_status": "Detection lab architecture finalized. AWS isolated VPC with 4 VMs (DC, 2x Win11, Kali), dedicated dev Elastic Cloud 9.3.2 cluster, development CrowdStrike CID. ART validation loops begin upon lab VM delivery.",
    "next_milestone": "DH-001 T1003.001 LSASS credential dump: first complete ART validation loop (POC for entire pipeline)"
  }
}